# supervisord.conf
[program:node]
command = node boot.js
directory = /home/dotcloud/current